# com.localbrain.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localbrain.api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/.local/share/local-ai-brain-prod &amp;&amp; export PATH="$HOME/.local/bin:$HOME/.cargo/bin:$PATH" &amp;&amp; exec uv run local-brain serve</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <!-- WARNING: macOS launchd does not rotate these log files automatically.
         Monitor file sizes or configure newsyslog(8) for rotation.
         The application also uses loguru with built-in rotation. -->
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/com.localbrain.api.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/com.localbrain.api.err.log</string>
</dict>
</plist>
